# docs/systemd/poly-trader.service
[Unit]
Description=Polymarket BTC 5-min Trader (v1.14)
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/poly
EnvironmentFile=/home/ubuntu/poly/.env
# Clear any stale lock from a previous OOM-kill / crash. The trader writes
# its lock as hostname:pid with no TTL, so a non-graceful exit leaves the
# key behind and the next start would fail with "lock held".
ExecStartPre=/usr/bin/redis-cli DEL poly:prod:trader:lock
ExecStart=/home/ubuntu/poly/target/release/poly-trader \
  --direction up \
  --base 10 \
  --max-step 5 \
  --band-min 0.25 \
  --band-max 0.75 \
  --exit-rule tp-sl \
  --tp-price 0.83 \
  --rsi-filter \
  --maker \
  --fixed-stake
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/poly/logs/poly-trader.systemd.log
StandardError=append:/home/ubuntu/poly/logs/poly-trader.systemd.log

# Resource limits — 1 GB RAM box, leave headroom for OS + redis
MemoryMax=512M
MemoryHigh=384M
CPUQuota=150%

[Install]
WantedBy=multi-user.target
